# zello.service
[Unit]
Description=Zello Bridge Service
[Service]
ExecStart=/path/to/zello_launcher.sh
[Install]
WantedBy=multi-user.target